# fiscalcidadao/settings/supervisor-program.cfg
[program:gunicorn]
environment=DJANGO_SETTINGS_MODULE=fiscalcidadao.settings.fiscalcidadao,PATH="/home/fiscalcidadao/Envs/fiscalcidadao/bin/"
command=/home/fiscalcidadao/Envs/fiscalcidadao/bin/gunicorn fiscalcidadao.wsgi:application -c /home/fiscalcidadao/fiscalcidadao-web/fiscalcidadao/settings/gunicorn.py
directory=/home/fiscalcidadao/fiscalcidadao-web/
user=fiscalcidadao
autostart=true
autorestart=true
redirect_stderr=True